# systemd/buryadvoice-bot.service
[Unit]
Description=Buryad Voice Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/BuryadVoiceBot
EnvironmentFile=/root/BuryadVoiceBot/.env
ExecStart=/root/BuryadVoiceBot/venv/bin/python bot.py
Restart=always
RestartSec=5
StandardOutput=append:/var/log/buryadvoice-bot.log
StandardError=append:/var/log/buryadvoice-bot-error.log

[Install]
WantedBy=multi-user.target